# contrib/init/bitcoind.conf
description "Bitcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env DETKCOIND_BIN="/usr/bin/detkcoind"
env DETKCOIND_USER="detkcoin"
env DETKCOIND_GROUP="detkcoin"
env DETKCOIND_PIDDIR="/var/run/detkcoind"
# upstart can't handle variables constructed with other variables
env DETKCOIND_PIDFILE="/var/run/detkcoind/detkcoind.pid"
env DETKCOIND_CONFIGFILE="/etc/detkcoin/detkcoin.conf"
env DETKCOIND_DATADIR="/var/lib/detkcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # detkcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$DETKCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run detkcoind."
        echo "The setting must appear in $DETKCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$DETKCOIND_PIDDIR"
    chmod 0755 "$DETKCOIND_PIDDIR"
    chown $DETKCOIND_USER:$DETKCOIND_GROUP "$DETKCOIND_PIDDIR"
    chown $DETKCOIND_USER:$DETKCOIND_GROUP "$DETKCOIND_CONFIGFILE"
    chmod 0660 "$DETKCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$DETKCOIND_PIDFILE" \
    --chuid $DETKCOIND_USER:$DETKCOIND_GROUP \
    --exec "$DETKCOIND_BIN" \
    -- \
    -pid="$DETKCOIND_PIDFILE" \
    -conf="$DETKCOIND_CONFIGFILE" \
    -datadir="$DETKCOIND_DATADIR" \
    -disablewallet \
    -daemon

